# deploy/legendas-backend.service
[Unit]
Description=Legendas Locais Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/legendas-locais/backend
EnvironmentFile=/opt/legendas-locais/backend/.env
ExecStart=/opt/legendas-locais/backend/.venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
